# conf/conf.supervisord.nginx.conf
[program:nginx]
command=/www/lanmps/action/nginx start
autostart=true
autorestart=true
stderr_logfile=/www/wwwLogs/supervisord_nginx_err.log
stdout_logfile=/www/wwwLogs/supervisord_nginx_out.log
priority=10
stdout_events_enabled=true
stderr_events_enabled=true